# debian/prerm
#!/bin/sh

if [ -d /run/systemd/system ] && [ "$1" = remove ]; then
   systemctl stop DCore || true
   systemctl disable DCore || true
fi
